# service/ospy.service
[Unit]
Description=OSPy irrigation controller
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory={{OSPY_DIR}}
ExecStart={{PYTHON}} -u {{OSPY_DIR}}/run.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=45
KillMode=control-group
SyslogIdentifier=ospy

[Install]
WantedBy=multi-user.target
